# misagodocker/crontab
0 0 * * * cd /misago && ./cron.sh >> /misago/cron.log 2>&1
